# battery_threshold/battery_thresholds.service
[Unit]
Description=Set battery charge thresholds
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/set_battery_thresholds.sh

[Install]
WantedBy=multi-user.target

